# ansible/templates/webapp.service.j2
[Unit]
Description={{ app_name }} Web Application
Documentation=https://github.com/{{ app_git_repo | basename }}
After=docker.service postgresql.service redis.service
Requires=docker.service
Wants=postgresql.service redis.service

[Service]
Type=simple
RemainAfterExit=no
Restart=always
RestartSec=10

# 停止旧容器
ExecStartPre=-/usr/bin/docker stop {{ app_name }}-app
ExecStartPre=-/usr/bin/docker rm {{ app_name }}-app

# 启动新容器
ExecStart=/usr/bin/docker run \
    --rm \
    --name {{ app_name }}-app \
    --network host \
    --env-file {{ app_directory }}/.env \
    -v {{ app_directory }}/data:/app/data \
    -v {{ app_directory }}/public:/app/public \
    -v {{ app_directory }}/uploads:/app/uploads \
    -p 127.0.0.1:{{ app_port }}:{{ app_port }} \
    {{ app_name }}:{{ app_git_branch | default('latest') }}

# 停止容器
ExecStop=/usr/bin/docker stop {{ app_name }}-app

# 健康检查
ExecStartPost=/bin/bash -c 'sleep 10 && curl -f http://localhost:{{ app_port }}/health || exit 1'

# 运行用户
User=deploy
Group=deploy

# 环境变量
Environment=HOME={{ app_directory }}

# 安全设置
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
